# installer/data/upstart/conpot.conf
########################################################
# T-Pot                                                #
# ConPot upstart script                                #
#                                                      #
# v16.03.2 by mo, DTAG, 2016-03-02                     #
########################################################

description "ConPot"
author "mo"
start on started docker and filesystem
stop on runlevel [!2345]
respawn
pre-start script
  # Remove any existing conpot containers
  myCID=$(docker ps -a | grep conpot | awk '{ print $1 }')
  if [ "$myCID" != "" ];
    then docker rm -v $myCID;
  fi
  # Remove any data from previous container if persistence is not enabled 
  if ! [ -f /data/persistence.on ];
    then
      rm -rf /data/conpot/* || true
      mkdir -p /data/conpot/log
      chmod 760 /data/conpot -R
      chown tpot:tpot /data/conpot -R
  fi
end script
script
  /usr/bin/docker run --name conpot --rm=true -v /data/conpot:/data/conpot -v /data/ews:/data/ews -p 81:80 -p 102:102 -p 161:161/udp -p 502:502 dtagdevsec/conpot:latest1603
end script
post-start script
  # Delay next start to avoid rapid respawning
  sleep 2
end script
