# host/etc/rc.local
#!/bin/bash
/opt/tpot/bin/updateip.sh
exit 0
