# installer/upstart/kippo.conf
########################################################
# T-Pot Community Edition                              #
# Kippo upstart script                                 #
#                                                      #
# v0.50 by mo, DTAG, 2015-01-27                        #
########################################################

description "Kippo"
author "mo"
start on started docker and filesystem
stop on runlevel [!2345]
respawn
pre-start script
  # Remove any existing kippo containers
  myCID=$(docker ps -a | grep kippo | awk '{ print $1 }')
  if [ "$myCID" != "" ];
    then docker rm -v $myCID;
  fi
end script
script
  # Delayed start to avoid rapid respawning
  sleep $(((RANDOM % 5)+5))
  /usr/bin/docker run --name kippo --rm=true -p 22:2222 -v /data/kippo dtagdevsec/kippo
end script
